# etc/cron.d/converter-cleanup
* * * * * root [ -z "$CONVERTER_ROOT" ] || find $CONVERTER_ROOT/data -type d -mmin +15 -exec rm -rf {} \; 

