# etc/cron.d/converter-cleanup
* * * * * find /var/www/md-ftw.info/data -mindepth 1 -type d -mmin +15 -exec rm -rf '{}' '+' 

